# systemd/solartracer-mqtt.service
[Unit]
Description=SolarTracer MQTT Publisher
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
User=@INSTALL_USER@
Group=@INSTALL_USER@
WorkingDirectory=@REPO_DIR@
ExecStart=/usr/bin/python3 @REPO_DIR@/publish_mqtt.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SupplementaryGroups=dialout

[Install]
WantedBy=multi-user.target